# artempd@.service
[Unit]
Description=Arduino temperature display service for %I

[Service]
Environment="PATH=/bin/"
ExecStart=/bin/artempd /dev/%I
Type=simple
